# firstboot/benchmark-firstboot.service
[Unit]
Description=Benchmark Thesis First-Boot Setup
After=network-online.target
Wants=network-online.target
ConditionPathExists=!/var/lib/benchmark-firstboot-complete

[Service]
Type=oneshot
ExecStart=/usr/local/bin/first-boot-setup.sh
TimeoutStartSec=600
User=root
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
